# embedded_browser/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start NGINX service
# ==============================================================================

# Wait for embedded browser to start before continuing
bashio::net.wait_for 5901

bashio::log.info "Starting Nginx..."
exec nginx
