# embedded_browser/rootfs/etc/services.d/websockify/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start websockify — WebSocket bridge between noVNC and TigerVNC
# ==============================================================================

# Wait until VNC server is up and listening on port 5900
bashio::log.info "Websockify waiting for VNC to start"
bashio::net.wait_for 5900

bashio::log.info "Starting websockify..."
exec websockify -v --web /usr/share/novnc/ 127.0.0.1:5901 127.0.0.1:5900
